# Imaging system/biolum_controller.service
[Unit]
Description=Biolum Camera Controller
After=network.target

[Service]
User=biolum
WorkingDirectory=/home/biolum
ExecStart=/usr/bin/python3 /home/biolum/biolum_controller.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
